# scripts/env-setup/crontab
*/5 * * * * root /root/swiftupload >> swift.log
