# scripts/env-setup/devd/docker/crontab
*/5 * * * * /palsweb/scripts/env-setup/devd/swiftupload >> swift.log
